# tectonic/services/packetbeat/vxlan.service
[Unit]
After=network.target

[Service]
ExecStart=sh /usr/local/bin/vxlan.sh

[Install]
WantedBy=default.target